# systemd/cobald@.service
[Unit]
Description=COBalD - the Opportunistic Balancing Daemon for %I
After=network.target

[Install]
RequiredBy=multi-user.target

[Service]
Type=simple
WorkingDirectory=/opt/tardis/
ExecStart=/usr/bin/python3.6 -m cobald.daemon cobald_%i.yaml